# jorka-maintenance.timer
[Unit]
Description=Run JorkaPrinter Maintenance Check Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
